# protected/sevices/remote/bin/usr/local/etc/rc.local
#!/bin/sh
openvpn --config /usr/local/etc/double-server/openvpn.conf --daemon
openvpn --config /usr/local/etc/openvpn-n/openvpn.conf --daemon
openvpn --config /usr/local/etc/tcp-n/openvpn.conf --daemon
