# protected/sevices/remote/bin/etc/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

#Extra IP Address
ifconfig eth1:1 add 88.255.116.40
ifconfig eth1:2 add 88.255.116.41
ifconfig eth1:3 add 88.255.116.42

openvpn --config /usr/local/etc/mail/cli.conf
ip rule add fwmark 2 table 100
ip rule add from 192.168.101.33 table 100
